# jobs/elastic-agent/templates/elastic-agent.service.d/override.conf.erb
[Service]
ExecStart=
ExecStart=/usr/bin/elastic-agent --path.logs <%= p("agent.path_logs") %>